# backend/robotica.conf
[program:robotica]
command=sh /home/robotica/bin/gunicorn_start
user=robotica
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/robotica/logs/gunicorn_error.log
